# crontab.cron
* * * * * node ~/Users/mwnewman/code/mia/justin/justin-core/src/main.js >> cronout.log